# EC_VOB/project_build_environment/aos5810_54x/fs/rootfs_addon/etc/systemd/system/serial-getty@ttyS1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/etc/init.d/rc.accton.init_deb
StandardInput=tty
StandardOutput=tty
Type=oneshot
Restart=no
RemainAfterExit=yes
