# EC_VOB/project_build_environment/aos5810_54x/fs/rootfs_addon/lib/systemd/system/platform-baseconf.service
[Unit]
Description=platform-baseconf Service
Before=serial-getty@ttyS1.service

[Service]
Type=oneshot
ExecStart=/etc/init.d/platform-baseconf
ExecStartPost=/bin/sleep 3
RemainAfterExit=yes
SuccessExitStatus=SIGHUP
StandardOutput=journal+console

[Install]
WantedBy=getty.target
